# ops/systemd/openai-crawler.service
[Unit]
Description=Paid Politely OpenAI Reddit Crawler
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/openai-crawler
Environment=NODE_ENV=production
EnvironmentFile=-/opt/openai-crawler/.env
ExecStart=/usr/bin/env bash -lc 'corepack pnpm crawl'
Restart=always
RestartSec=10
KillSignal=SIGTERM
TimeoutStopSec=30
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
